# schnapps/schnapps.cron
MAILTO=""
# m h    dom mon dow  user  command
  * */12  *   *   *   root  schnapps cleanup
  0 1     *   *   0   root  schnapps create -t time "Snapshot created by cron"
