# watchdog/watchdog.cron
*/5     *       *       *       *       /usr/bin/watchdog.sh >/dev/null 2>&1
