# watchdog/watchdog.cron
*/5     *       *       *       *       root	/usr/bin/watchdog.sh >/dev/null 2>&1
