# logsend/logsend.cron
MAILTO=michal.vaner@nic.cz
5     *       *       *       *       root	/usr/bin/logsend.sh
