# logsend/logsend.cron
MAILTO=""
20     *       *       *       *       root	/usr/bin/logsend.sh
