# logsend/logsend.cron
MAILTO=""
5     *       *       *       *       root	/usr/bin/logsend.sh
